# templates/default/init.sh.erb
#! /bin/sh

### BEGIN INIT INFO
# Provides:             sshd
# Required-Start:       $remote_fs $syslog
# Required-Stop:        $remote_fs $syslog
# Default-Start:        2 3 4 5
# Default-Stop:         
# Short-Description:    SQLAnywhere server
### END INIT INFO

set -e

# /etc/init.d/ssh: start and stop the SQLAnywhere server daemon
SA12_INSTALL="<%= node['sqlanywhere']['install_dir'] %>"

test -x $SA12_INSTALL/<%= node['sqlanywhere']['bin_dir'] %>/dbsrv12 || exit 0
exec /lib/init/upstart-job sa12 "$@"